# scripts/calendarspread.service
[Unit]
Description=CalendarSpread diff-z reversion live trading loop (autonomous)
After=network-online.target
Wants=network-online.target
# Crash-loop guard MUST live in [Unit] — in [Service] systemd ignores it
# (that's the "Unknown key name 'StartLimitIntervalSec'" warning in journald).
StartLimitIntervalSec=600
StartLimitBurst=5

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/CalendarSpread
# Create the log dir before any file/journal handle opens (fixes the past
# 209/STDOUT crash when logs/ was missing on a fresh boot).
ExecStartPre=/usr/bin/mkdir -p /home/ec2-user/CalendarSpread/logs
# Autonomous diff-z reversion with P2 accumulation. --size-from-wallet sizes off
# live free USDC × --wallet-frac, which sets the PER-EVENT worst-case-loss cap
# (5% of free wallet per event). Ratchet --wallet-frac up as closed_pnl confirms.
ExecStart=/home/ec2-user/CalendarSpread/.venv/bin/python3 live_execution.py --execute --loop-seconds 600 --size-from-wallet --wallet-frac 0.05 --accumulate

# A trading loop should never stay down — restart on ANY exit.
Restart=always
RestartSec=30

# Output to journald (view with `journalctl -u calendarspread -f`)
StandardOutput=journal
StandardError=journal

# Resource limits
MemoryMax=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
